# tools/cron.d/mysql_backup.cron
14 3 * * * project cd src && tools/cron.scripts/mysql_backup.sh 2>&1 | logger -t 'mysql_backup.cron'

